# scripts/launchd/com.foundrysmb.psygil.periodic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.foundrysmb.psygil.periodic</string>
	<key>ProgramArguments</key>
	<array>
		<string>__REPO_ROOT__/scripts/psygil-sync.sh</string>
		<string>commit</string>
	</array>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>StandardOutPath</key>
	<string>/tmp/psygil-periodic.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/psygil-periodic-err.log</string>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
