# scripts/launchd/com.foundrysmb.psygil.sleep-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.foundrysmb.psygil.sleep-push</string>
	<key>ProgramArguments</key>
	<array>
		<string>__REPO_ROOT__/scripts/psygil-sync.sh</string>
		<string>push</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/Library/Preferences/com.apple.powermanagement.plist</string>
	</array>
	<key>StandardOutPath</key>
	<string>/tmp/psygil-sleep-push.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/psygil-sleep-push-err.log</string>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
